# deploy/systemd/tmuxbot.service
# tmuxbot — systemd user service
#
# 安装:
#   mkdir -p ~/.config/systemd/user
#   ln -s ABS_PATH_TO_REPO/deploy/systemd/tmuxbot.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now tmuxbot.service
#   loginctl enable-linger $USER   # 让 service 在 logout 后继续跑
#
# 看 log:  journalctl --user -u tmuxbot -f
# 重启:    systemctl --user restart tmuxbot
# 停:      systemctl --user stop tmuxbot
#
# 推荐通过 `tmuxbot install-service --now` 安装；本文件可手工复制。

[Unit]
Description=tmuxbot WebUI and tmux bridge supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Keep native tmux/CLI sessions alive when the bridge service is restarted.
KillMode=process
EnvironmentFile=-%h/.config/tmuxbot/.env
Environment=TMUXBOT_BRIDGE_PID_FILE=%t/tmuxbot/bridge.pid
ExecStop=/bin/sh -c 'if [ -r "$TMUXBOT_BRIDGE_PID_FILE" ]; then kill -TERM "$(cat "$TMUXBOT_BRIDGE_PID_FILE")" 2>/dev/null || true; fi'
ExecStart=%h/.local/bin/tmuxbot serve
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
# 资源上限 — 多 backend / 多 binding 后可适当调高
MemoryHigh=2G
MemoryMax=4G
# 服务退出时 tmux session 保持存活。

[Install]
WantedBy=default.target
