# deploy/systemd/tmuxbot-web.service
# tmuxbot Web control plane — systemd user service
#
# 安装:
#   mkdir -p ~/.config/systemd/user
#   ln -s ABS_PATH_TO_REPO/deploy/systemd/tmuxbot-web.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now tmuxbot-web.service
#   loginctl enable-linger $USER   # 让 service 在 logout 后继续跑
#
# 看 log:  journalctl --user -u tmuxbot-web -f
# 重启:    systemctl --user restart tmuxbot-web
# 停:      systemctl --user stop tmuxbot-web
#
# 兼容模式：仅启动 Web；新部署请使用 tmuxbot.service。

[Unit]
Description=tmuxbot Web control plane
After=network.target

[Service]
Type=simple
KillMode=process
EnvironmentFile=-%h/.config/tmuxbot/.env
Environment=TMUXBOT_BRIDGE_PID_FILE=%t/tmuxbot/bridge.pid
ExecStop=/bin/sh -c 'if [ -r "$TMUXBOT_BRIDGE_PID_FILE" ]; then kill -TERM "$(cat "$TMUXBOT_BRIDGE_PID_FILE")" 2>/dev/null || true; fi'
ExecStart=%h/.local/bin/tmuxbot web
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
