# deploy/systemd/tmuxbot-codex-rotation.service
[Unit]
Description=Refresh tmuxbot Codex Feishu WebSocket bridge

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart tmuxbot-codex.service
