# deploy/systemd/tmuxbot-bridge-refresh@.service
[Unit]
Description=Refresh tmuxbot channel bridge %i

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart %i.service
