# deploy/systemd/tmuxbot-bridge-refresh@.timer
[Unit]
Description=Refresh tmuxbot channel bridge %i every six hours

[Timer]
OnBootSec=30min
OnUnitInactiveSec=6h
Persistent=true
Unit=tmuxbot-bridge-refresh@%i.service

[Install]
WantedBy=timers.target
